# deploy/e2i-bentoml.service
[Unit]
Description=E2I BentoML Model Serving Service
Documentation=https://github.com/e2i/e2i_causal_analytics
After=network.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=enunez
Group=enunez
WorkingDirectory=/opt/e2i_causal_analytics
EnvironmentFile=/opt/e2i_causal_analytics/deploy/e2i-bentoml.env
Environment="PATH=/opt/e2i_causal_analytics/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

ExecStart=/opt/e2i_causal_analytics/.venv/bin/bentoml serve \
    e2i_serving_service.py:E2IModelService \
    --host 0.0.0.0 \
    --port 3000 \
    --working-dir /opt/e2i_causal_analytics/scripts/bentoml

Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=e2i-bentoml

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/enunez/bentoml /opt/e2i_causal_analytics /tmp
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
